# scripts/com.foval.gtd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.foval.gtd</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/johnfoval/Desktop/Nurik/gtd-app/node_modules/.bin/next</string>
    <string>start</string>
    <string>-H</string>
    <string>0.0.0.0</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johnfoval/Desktop/Nurik/gtd-app</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/johnfoval/Desktop/Nurik/gtd-app/data/logs/gtd-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnfoval/Desktop/Nurik/gtd-app/data/logs/gtd-stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
</dict>
</plist>
